# files/etc/init.d/blue-merle
#!/bin/sh /etc/rc.common

. /lib/blue-merle/functions.sh

# We intend to be started before the first network-related service is started.
# According to https://openwrt.org/docs/techref/initscripts, /etc/rc.d/ determines
# the order of the services to be started (or stopped). The lower the number,
# the earlier the service is started.
# We observe "repeater" having the value 15. "network" 20. We certainly want to ahead of those.
START=10
STOP=99

# Report the SSID that this boot picked on the MCU screen so the user
# knows what to look for in their WiFi picker. Best-effort: skipped if
# /dev/ttyS0 is not a character device (e.g. not a real Mudi).
_announce_ssid_on_mcu() {
	[ -c /dev/ttyS0 ] || return 0
	local ssid
	ssid=$(uci -q get wireless.@wifi-iface[0].ssid) || return 0
	# MCU display is 16 x 2 characters; the JSON wrapper trims long
	# strings automatically. Format: label on line 1, SSID on line 2.
	printf '{"msg":"WiFi:           %s"}\n' "$ssid" > /dev/ttyS0
}

# Read config main.main.stable_identity (default 0). When 1, boot-time
# rotation is skipped for EVERY identifier — hostname, SSID, BSSIDs
# AND client-visible MACs — so the user gets one truly stable Apple-
# style identity until they manually run `blue-merle-newmac --full`.
#
# An earlier version still called RANDOMIZE_MACADDR under
# stable_identity=1 on the theory that client-side MACs don't leak
# across sessions like SSID does. That was inconsistent with the
# "stable identity" contract users would reasonably expect (a real
# iPhone doesn't rotate its client MAC every boot either), so we now
# leave MACs alone too. Users who want partial mixing can call
# blue-merle-newmac (without --full) manually.
#
# See README.md § "Threat-model tradeoff" for the rationale.
_stable_identity_enabled() {
	local v
	v=$(uci -q get blue-merle.main.stable_identity 2>/dev/null || echo 0)
	[ "$v" = "1" ]
}

start() {
	if _stable_identity_enabled; then
		logger -p notice -t blue-merle \
			"stable_identity=1: no rotation on boot; use blue-merle-newmac --full to rotate manually"
	else
		RESET_BSSIDS
		RANDOMIZE_MACADDR
		RANDOMIZE_HOSTNAME
		RANDOMIZE_SSID
	fi
	_announce_ssid_on_mcu
}

# `service blue-merle reload` (or the hotplug hook below) regenerates the
# BSSIDs/MACs without a full reboot. Useful when the user restarts WiFi
# manually and doesn't want the previous BSSID to linger.
reload() {
	if _stable_identity_enabled; then
		logger -p notice -t blue-merle \
			"stable_identity=1: reload is a no-op; use blue-merle-newmac --full to rotate manually"
	else
		RESET_BSSIDS
		RANDOMIZE_MACADDR
		RANDOMIZE_SSID
	fi
}
